# extras/stats.cron
*/5 * * * * root  /opt/quicklook/stats.py --data=/var/opt/quicklook/ --output=/var/www/quicklook/ > /dev/null
